# src/config/autostart/clj_de_autostart.sh.desktop
[Desktop Entry]
Exec=clj_de_autostart
Icon=dialog-scripts
Name=clj_de_autostart
Path=
Type=Application
X-KDE-AutostartScript=true
